# supervisor/conf/development.conf
[program:fast-laravel-monitor]
command=/usr/local/bin/php fast:http start ; the program (relative uses PATH, can take args)
;process_name=%(program_name)s  ; process_name expr (default %(program_name)s)
numprocs=1                     ; number of processes copies to start (def 1)
directory=/Users/xmc/PhpstormProjects/test/api                 ; directory to cwd to before exec (def no cwd)
;umask=022                      ; umask for process (default None)
;priority=999                   ; the relative start priority (default 999)
;exitcodes=0,2                  ; 'expected' exit codes used with autorestart (default 0,2)
stopsignal=SIGTERM                ; signal used to kill process (default TERM)
;stopwaitsecs=10                ; max num secs to wait b4 SIGKILL (default 10)
;stopasgroup=false              ; send stop signal to the UNIX process group (default false)
;killasgroup=false              ; SIGKILL the UNIX process group (def false)
;stdout_capture_maxbytes=1MB    ; number of bytes in 'capturemode' (default 0)
;stdout_events_enabled=false    ; emit events on stdout writes (default false)
;stderr_logfile=NONE            ; stderr log path, NONE for none; default AUTO
;stderr_logfile_maxbytes=1MB    ; max # logfile bytes b4 rotation (default 50MB)
;stderr_logfile_backups=10      ; # of stderr logfile backups (default 10)
;stderr_capture_maxbytes=1MB    ; number of bytes in 'capturemode' (default 0)
;stderr_events_enabled=false    ; emit events on stderr writes (default false)
;environment=A="1",B="2"        ; process environment additions (def no adds)
;serverurl=AUTO                 ; override serverurl computation (childutils)

autostart=true                  ; start at supervisord start (default: true)
startsecs=1                     ; # of secs prog must stay up to be running (def. 1)
startretries=3                  ; max # of serial start failures when starting (default 3)
autorestart=unexpected          ; when to restart if exited after running (def: unexpected)
user=xmc                       ; setuid to this UNIX account to run the program
redirect_stderr=true            ; redirect proc stderr to stdout (default false)
stdout_logfile= /Users/xmc/PhpstormProjects/test/api/storage/logs/fast-laravel-monitor.log      ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10       ; # of stdout logfile backups (default 10)